# systemd/.config/systemd/user/kanshi.service
[Unit]
Description=Kanshi display manager
After=graphical-session.target

[Service]
ExecStart=/usr/local/bin/kanshi
Restart=always
RestartSec=5
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/%U

[Install]
WantedBy=graphical-session.target
